# root/etc/services.d/pkg-cache-error-log/run
#!/bin/sh
exec tail -F /config/log/error.log >&2
